# aur-updater/aur-updater.timer
[Unit]
Description=Daily 20:00 AUR auto-updater run
Documentation=file://%h/work/github.com/Xefreh/pkgbuilds/aur-updater/README.md

[Timer]
# Every day at 20:00. Persistent=true catches up a run missed while the machine
# was asleep/off (it runs at next boot).
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
