# aur-updater/aur-updater.service
[Unit]
Description=AUR auto-updater (pkgver bump + push AUR + Telegram digest)
Documentation=file://%h/work/github.com/Xefreh/pkgbuilds/aur-updater/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Repo root containing the package submodules.
WorkingDirectory=%h/work/github.com/Xefreh/pkgbuilds
# Telegram secrets read from the env (never committed).
EnvironmentFile=%h/.config/aur-updater.env
ExecStart=%h/work/github.com/Xefreh/pkgbuilds/aur-updater/target/release/aur-updater
# A failed run does not fail the unit; the timer keeps going.
SuccessExitStatus=0
